# I2C/OLED_R2S
#!/bin/sh /etc/rc.common

START=99
start() {
	while :; do
		nohup /usr/bin/i2c_ssd 1>/dev/null 2>&1 &
		sleep 15m
		killall -q -9 i2c_ssd
	done
}
